# apps/activemq/cluster/mq_sync.service
[Unit]
Description=MQ Sync Service

[Service]
ExecStart=/path/to/mq_sync.sh
