# apps/activemq/cluster/rclone-sync.timer
[Unit]
Description=Runs mq_sync every second

[Timer]
OnUnitActiveSec=1s
Unit=mq_sync.service

[Install]
WantedBy=timers.target
