# apps/activemq/cluster/rclone-sync.service
[Unit]
Description=Rclone Service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart={path}/mq_sync.sh
Restart=on-failure
RestartSec=3
User=root
Group=root

[Install]
WantedBy=multi-user.target
